# systemd/super/journalctl-log.timer
[Unit]
Description=pipe journalctl output to file once every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
