# rtorrent/torrent.service
[Unit]
Description=starts torrent clients

[Service]
ExecStart=/home/aidan/torrent/rtorrent-service.sh

[Install]
WantedBy=multi-user.target
